# scripts/launchd/com.marin-monitor.sync-activity.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.marin-monitor.sync-activity</string>

	<key>ProgramArguments</key>
	<array>
		<string>/Users/tammypais/projects/marin-monitor/scripts/sync-runner.sh</string>
		<string>activity</string>
	</array>

	<key>StartInterval</key>
	<integer>21600</integer>

	<key>WorkingDirectory</key>
	<string>/Users/tammypais/projects/marin-monitor</string>

	<key>StandardOutPath</key>
	<string>/Users/tammypais/Library/Logs/marin-monitor/sync-activity-launchd.log</string>

	<key>StandardErrorPath</key>
	<string>/Users/tammypais/Library/Logs/marin-monitor/sync-activity-launchd.log</string>

	<key>RunAtLoad</key>
	<true/>
</dict>
</plist>
